# com.nse.screener.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nse.screener</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.13/bin/python3</string>
        <string>-m</string>
        <string>uvicorn</string>
        <string>app.main:app</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>9000</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ajayagarwal/Desktop/Stocks project/nse-fno-screener</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPYCACHEPREFIX</key>
        <string>/tmp/nse-fno-pycache</string>
        <key>HOME</key>
        <string>/Users/ajayagarwal</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/nse-screener.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/nse-screener.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
